# overlay.d/15fcos/usr/lib/udev/rules.d/65-gdm.rules
# Re-enable watland for nvidia (run after 61-gdm.rules)
DRIVER=="nvidia", RUN+="/usr/libexec/gdm-runtime-config set daemon WaylandEnable true"