# builds/client/resources/62-gdm-nvidia-wayland.rules
# re-enable Wayland when using the proprietary nvidia driver
DRIVER=="nvidia", RUN+="/usr/libexec/gdm-runtime-config set daemon WaylandEnable true"
