# overlay.d/01custom/usr/lib/dracut/modules.d/35coreos-network/coreos-copy-firstboot-network.service
[Unit]
Description=Copy CoreOS Firstboot Networking Config
ConditionPathExists=/usr/lib/initrd-release
DefaultDependencies=false

[Service]
Type=oneshot
RemainAfterExit=yes
MountFlags=slave
ExecStart=/usr/sbin/coreos-copy-firstboot-network
