# overlay.d/01custom/usr/lib/dracut/modules.d/35coreos-network/coreos-enable-network.service
[Unit]
Description=CoreOS Enable Network
ConditionPathExists=/etc/initrd-release
DefaultDependencies=false

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/coreos-enable-network
